# rules/registerjarinrunkey.yml
title: Register Jar In Run Key
status: experimental
description: register jar in run key
author: Joe Security
date: 2019-11-07
id: 200022
threatname:
behaviorgroup: 2
classification: 4
mitreattack: T1547.001

logsource:
    category: process_creation
    product: windows
detection:
    selection:
        CommandLine:
            - '*reg add *\Software\Microsoft\Windows\CurrentVersion\Run /v *\AppData\Roaming\\*.jar*'
            - '*reg add *\software\microsoft\windows\currentversion\run /v *javaw.exe -jar*'
    condition: selection
level: critical
